# openwrt/ulldpd/files/etc/hotplug.d/iface/20-lldp
#!/bin/sh

[ "$ACTION" = ifup -o "$ACTION" = ifupdate ] || exit 0
ubus call lldp check_devices
